# provision/ansible/docker/roles/os/templates/nfs/mnt-media.mount.j2
[Unit]
Description=Mount media at boot

[Mount]
What={{ nfs_src }}
Where={{ nfs_mountpoint }}
Options=vers=4
Type=nfs
TimeoutSec=30

[Install]
WantedBy=multi-user.target